# crontab
* * * * 1-5 cp /etc/hostsoff /etc/hosts

* * * * 6-7 cp /etc/hostson /etc/host
